# cec-initialize.service
[Unit]
Description=Initialize CEC device

[Service]
Type=oneshot
ExecStart=/usr/bin/cec-ctl --playback
RemainAfterExit=yes
