# cec-on.timer
[Unit]
Description=Turn on CEC controlled device 

[Timer]
OnCalendar=Mon,Tue,Wed,Thu,Fri *-*-* 7:00:00

[Install]
WantedBy=timers.target