# cec-off.service
[Unit]
Description=Turns off cec controlled TV 

[Service]
Type=oneshot
ExecStart=/usr/bin/cec-ctl --standby -t0 
